# roles/ssh_server/files/ssh.service.d/override.conf
[Unit]
After=network-online.target auditd.service
